# config/labwc-screen.service
[Unit]
Description=labwc Wayland compositor (monitor físico)
After=graphical.target

[Service]
Type=simple
Environment=XDG_SESSION_TYPE=wayland
ExecStart=/usr/bin/labwc
Restart=on-failure
RestartSec=5
